# src/supervisord-smbd.conf
[program:smbd]
command=/etc/init.d/smbd start
autostart=true
autorestart=true

[program:nmbd]
command=/etc/init.d/nmbd start
autostart=true
autorestart=true
